# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/dockeruser/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Web terminal (ttyd) ──────────────────────────────────────────────────────

[program:ttyd]
command=ttyd -i 0.0.0.0 -p 8080 -T xterm-256color -W bash
autostart=true
autorestart=true
priority=50
stdout_logfile=/home/dockeruser/logs/ttyd.log
redirect_stderr=true

; ── Virtual display (must start first) ───────────────────────────────────────

[program:xvfb]
command=Xvfb :98 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/home/dockeruser/logs/xvfb.log
redirect_stderr=true

; ── Window manager on virtual display ────────────────────────────────────────

[program:openbox]
command=openbox --sm-disable
environment=DISPLAY=":98",NO_AT_BRIDGE="1"
autostart=true
autorestart=true
priority=20
stdout_logfile=/home/dockeruser/logs/openbox.log
redirect_stderr=true

; ── VNC server — captures virtual display ────────────────────────────────────

[program:x11vnc]
command=x11vnc -display :98 -forever -shared -nopw -rfbport 5901
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/home/dockeruser/logs/x11vnc.log
redirect_stderr=true

; ── WebSocket bridge — noVNC connects here ───────────────────────────────────

[program:websockify]
command=websockify --web=/opt/rviz-novnc 6090 localhost:5901
autostart=true
autorestart=true
priority=30
startretries=10
stdout_logfile=/home/dockeruser/logs/websockify.log
redirect_stderr=true
